# systemd/network_check.service
[Unit]
Description=Network Connection Monitor
After=network.target

[Service]
ExecStart=/usr/local/bin/network_check
Restart=always
User=root
Group=root

[Install]
WantedBy=multi-user.target

